# crontab
0 0 * * * /opt/backups/scripts/backup-windowspc.sh
0 2 * * * /opt/backups/scripts/backup-laptop.sh
0 18 * * 6 /opt/backups/scripts/purge-old-backups.sh
